# server-python/systemd/guildmarks-api.service
[Unit]
Description=GuildMarks dedicated API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=guildmarks
Group=guildmarks
WorkingDirectory=/opt/guildmarks-api
EnvironmentFile=/etc/guildmarks/guildmarks.env
ExecStart=/opt/guildmarks-api/.venv/bin/uvicorn guildmarks_api.app:app --host 127.0.0.1 --port 8765 --workers 1 --proxy-headers --forwarded-allow-ips 127.0.0.1 --no-access-log
Restart=on-failure
RestartSec=3
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
ReadWritePaths=/var/lib/guildmarks
ReadOnlyPaths=/etc/guildmarks /opt/guildmarks-api
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
